# docker/rc.local
#!/bin/bash

# fixes for elasticsearch
sysctl -w vm.max_map_count=262144

cd /root/elk-docker

docker-compose up -d


#containers (don't use , to separate items)
containers=(
  "docker-ui"

)

for name in "${containers}"
do
docker start $name
done

